# deploy/last-agent.service
# Run the agent as a native binary via systemd (no Docker).
#
#   sudo cp target/release/vantage-agent /usr/local/bin/
#   sudo cp deploy/vantage-agent.service /etc/systemd/system/
#   sudo mkdir -p /etc/vantage && sudoedit /etc/vantage/agent.env   # HUB_URL, API_KEY
#   sudo systemctl enable --now vantage-agent
#
[Unit]
Description=vantage agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/vantage/agent.env
ExecStart=/usr/local/bin/vantage-agent
Restart=always
RestartSec=5
# Hardening: the agent only needs to read /proc and reach the hub.
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
